# .provision/sharedfolders.conf
start on vagrant-mounted
pre-start script
  /bin/mount -t vboxsf -o uid=10000,gid=10000 home_trinity_src /home/trinity/src
end script
